# deploy/systemd/watchpuppy-onvif-gated.target
[Unit]
Description=WatchPuppy ONVIF gated collection target
