# deploy/systemd/watchpuppy-onvif-gated@.service
[Unit]
Description=WatchPuppy ONVIF gated collector for camera %i
Wants=network-online.target
Wants=watchpuppy-cnn-inference.service
After=network-online.target hailort.service watchpuppy-cnn-inference.service
PartOf=watchpuppy-onvif-gated.target

[Service]
Type=simple
User=moai
WorkingDirectory=/home/moai/Workspace/Codex/WatchPuppy
EnvironmentFile=/home/moai/Workspace/Codex/WatchDog/configs/watchdog-tapo.env
EnvironmentFile=/home/moai/Workspace/Codex/WatchPuppy/configs/watchpuppy.env
ExecStart=/usr/bin/env bash /home/moai/Workspace/Codex/WatchPuppy/scripts/run_onvif_watchpuppy_service.sh %i
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=10
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
